# systemd/clipboard-sync.service.template
[Unit]
Description=Clipboard sync between Linux and macOS
After=network.target

[Service]
Type=simple
ExecStart={{INSTALL_PREFIX}}/bin/clipboard-sync
Restart=on-failure
RestartSec=10
Environment=DISPLAY=:0
Environment=WAYLAND_DISPLAY=wayland-0
Environment=XDG_RUNTIME_DIR=/run/user/1000

[Install]
WantedBy=default.target
